# root/etc/s6-overlay/s6-rc.d/svc-esde/run
#!/usr/bin/with-contenv bash
set -euo pipefail

display_num="${DISPLAY_NUM:-0}"
export DISPLAY="${DISPLAY:-:${display_num}}"

if [[ -z ${XDG_RUNTIME_DIR:-} ]]; then
  echo "[svc-esde] Warning: XDG_RUNTIME_DIR not set. Setting to /run/esde" >&2
  export XDG_RUNTIME_DIR=/run/esde
fi

mkdir -p "${XDG_RUNTIME_DIR}"
chmod 700 "${XDG_RUNTIME_DIR}"
chown "${APP_USER:-abc}":"${APP_USER:-abc}" "${XDG_RUNTIME_DIR}"

# --- Device group membership ---
# The abc user needs supplementary groups for GPU, audio, and input access.
# s6-setuidgid reads supplementary groups from /etc/group, so the user must
# be a member before we drop privileges.
app_user="${APP_USER:-abc}"
for grp in video render audio input; do
  if getent group "${grp}" >/dev/null 2>&1; then
    if ! id -nG "${app_user}" 2>/dev/null | grep -qw "${grp}"; then
      usermod -aG "${grp}" "${app_user}" 2>/dev/null || true
    fi
  fi
done

# Match host device GIDs for /dev/dri and /dev/snd if they differ from standard groups
if [[ -e /dev/dri/card0 ]]; then
  dri_gid="$(stat -c '%g' /dev/dri/card0 2>/dev/null)" || true
  if [[ -n "${dri_gid}" && "${dri_gid}" != "0" ]]; then
    if ! id -G "${app_user}" 2>/dev/null | grep -qw "${dri_gid}"; then
      # Create/reuse a group matching the host GID and add the user
      dri_grp="$(getent group "${dri_gid}" 2>/dev/null | cut -d: -f1)" || true
      if [[ -z "${dri_grp}" ]]; then
        groupadd -g "${dri_gid}" hostdri 2>/dev/null || true
        dri_grp="hostdri"
      fi
      usermod -aG "${dri_grp}" "${app_user}" 2>/dev/null || true
      echo "[svc-esde] Added ${app_user} to group ${dri_grp} (GID ${dri_gid}) for /dev/dri access."
    fi
  fi
fi

if [[ -e /dev/snd/controlC0 ]]; then
  snd_gid="$(stat -c '%g' /dev/snd/controlC0 2>/dev/null)" || true
  if [[ -n "${snd_gid}" && "${snd_gid}" != "0" ]]; then
    if ! id -G "${app_user}" 2>/dev/null | grep -qw "${snd_gid}"; then
      snd_grp="$(getent group "${snd_gid}" 2>/dev/null | cut -d: -f1)" || true
      if [[ -z "${snd_grp}" ]]; then
        groupadd -g "${snd_gid}" hostaudio 2>/dev/null || true
        snd_grp="hostaudio"
      fi
      usermod -aG "${snd_grp}" "${app_user}" 2>/dev/null || true
      echo "[svc-esde] Added ${app_user} to group ${snd_grp} (GID ${snd_gid}) for /dev/snd access."
    fi
  fi
fi

# --- Audio detection and setup ---
# ES-DE requires PulseAudio or PipeWire; raw ALSA is not supported.

# Start udevd early so PulseAudio's module-udev-detect can discover ALSA cards,
# and Xorg can later discover input devices.
if ! pidof systemd-udevd >/dev/null 2>&1 && ! pidof udevd >/dev/null 2>&1; then
  if [[ -x /usr/lib/systemd/systemd-udevd ]]; then
    echo "[svc-esde] Starting udevd for device discovery."
    /usr/lib/systemd/systemd-udevd --daemon 2>/dev/null || true
    udevadm trigger --action=add 2>/dev/null || true
    udevadm settle --timeout=5 2>/dev/null || true
  elif command -v udevd >/dev/null 2>&1; then
    echo "[svc-esde] Starting udevd for device discovery."
    udevd --daemon 2>/dev/null || true
    udevadm trigger --action=add 2>/dev/null || true
    udevadm settle --timeout=5 2>/dev/null || true
  else
    echo "[svc-esde] Warning: udevd not found. Audio detection and input hotplug may not work." >&2
  fi
fi

# 1. Detect ALSA hardware (sound cards visible to the kernel)
if [[ -d /dev/snd ]]; then
  snd_devices="$(ls /dev/snd/ 2>/dev/null)" || true
  if [[ -n "${snd_devices}" ]]; then
    echo "[svc-esde] Audio hardware detected in /dev/snd: $(echo "${snd_devices}" | tr '\n' ' ')"
  else
    echo "[svc-esde] Warning: /dev/snd exists but contains no device nodes." >&2
  fi
  if [[ -r /proc/asound/cards ]]; then
    while IFS= read -r line; do
      echo "[svc-esde]   ALSA card: ${line}"
    done < /proc/asound/cards
  fi
else
  echo "[svc-esde] Warning: /dev/snd not found. No sound devices available." >&2
  echo "[svc-esde]   For USB audio: pass --device /dev/snd:/dev/snd or use privileged mode." >&2
fi

# 2. Start PulseAudio
if command -v pulseaudio >/dev/null 2>&1; then
  if ! s6-setuidgid "${APP_USER:-abc}" pulseaudio --check 2>/dev/null; then
    s6-setuidgid "${APP_USER:-abc}" pulseaudio --start --exit-idle-time=-1 --daemonize=yes 2>/dev/null || \
      echo "[svc-esde] Warning: PulseAudio failed to start. Audio may be unavailable." >&2
  fi

  # 3. Configure default sink (wait briefly for PulseAudio to discover devices)
  if s6-setuidgid "${APP_USER:-abc}" pulseaudio --check 2>/dev/null; then
    sleep 1

    # Log all discovered sinks for troubleshooting
    sink_list="$(s6-setuidgid "${APP_USER:-abc}" pactl list short sinks 2>/dev/null)" || true
    if [[ -n "${sink_list}" ]]; then
      echo "[svc-esde] PulseAudio sinks available:"
      while IFS= read -r sink_line; do
        echo "[svc-esde]   ${sink_line}"
      done <<< "${sink_list}"
    else
      echo "[svc-esde] Warning: PulseAudio running but no sinks detected." >&2
    fi

    # Set explicit sink if PULSE_SINK is defined, otherwise auto-select the first one
    if [[ -n "${PULSE_SINK:-}" ]]; then
      s6-setuidgid "${APP_USER:-abc}" pactl set-default-sink "${PULSE_SINK}" 2>/dev/null && \
        echo "[svc-esde] PulseAudio default sink set to: ${PULSE_SINK}" || \
        echo "[svc-esde] Warning: Could not set PULSE_SINK=${PULSE_SINK}. Check sink names above." >&2
    elif [[ -n "${sink_list}" ]]; then
      default_sink="$(echo "${sink_list}" | head -1 | awk '{print $2}')" || true
      if [[ -n "${default_sink:-}" ]]; then
        s6-setuidgid "${APP_USER:-abc}" pactl set-default-sink "${default_sink}" 2>/dev/null || true
        echo "[svc-esde] PulseAudio auto-selected sink: ${default_sink}"
      fi
    fi
  fi
else
  echo "[svc-esde] Warning: PulseAudio not installed. Audio will not work." >&2
fi

if [[ ! -e /dev/dri/renderD128 && ! -e /dev/dri/card0 ]]; then
  echo "[svc-esde] Warning: No /dev/dri render/card node detected. Hardware GL may be unavailable." >&2
fi

esde_home="${ESDE_HOME:-/config}"
read -r -a extra_args <<< "${ESDE_ARGS:-}"

# Ensure the ES-DE home directory is writable by the application user
# so ES-DE can create subdirectories such as ROMs, gamelists, etc.
if [[ -d "${esde_home}" ]]; then
  chown "${APP_USER:-abc}":"${APP_USER:-abc}" "${esde_home}" 2>/dev/null || true
fi

# --- Symlink external volume mounts into ES-DE's expected paths ---
# docker-compose mounts /roms and /bios as named volumes, but ES-DE expects
# ROMs under ${ESDE_HOME}/ROMs and BIOS files in emulator-specific paths.
# Create symlinks so both locations work seamlessly.
if [[ -d /roms && ! -e "${esde_home}/ROMs" ]]; then
  ln -s /roms "${esde_home}/ROMs"
  chown -h "${APP_USER:-abc}":"${APP_USER:-abc}" "${esde_home}/ROMs"
  echo "[svc-esde] Linked /roms -> ${esde_home}/ROMs"
fi
if [[ -d /bios && ! -e "${esde_home}/.emulationstation/bios" ]]; then
  mkdir -p "${esde_home}/.emulationstation"
  ln -s /bios "${esde_home}/.emulationstation/bios"
  chown -h "${APP_USER:-abc}":"${APP_USER:-abc}" "${esde_home}/.emulationstation" "${esde_home}/.emulationstation/bios"
  echo "[svc-esde] Linked /bios -> ${esde_home}/.emulationstation/bios"
fi

# --- Seed bundled demo ROMs on first boot ---
# /defaults/roms/ ships with the image and contains freely-licensed homebrew.
# Copy them into /roms (the volume) only if the volume is completely empty,
# so user-supplied content is never overwritten.
if [[ -d /defaults/roms && -d /roms ]]; then
  _rom_file_count="$(find /roms -mindepth 1 -maxdepth 2 -type f \( -name '*.gb' -o -name '*.gbc' -o -name '*.nes' -o -name '*.sfc' -o -name '*.smc' -o -name '*.bin' -o -name '*.zip' -o -name '*.7z' \) 2>/dev/null | wc -l)"
  if [[ "${_rom_file_count}" -eq 0 ]]; then
    echo "[svc-esde] No ROM files found in /roms; seeding bundled demo ROMs..."
    cp -rn /defaults/roms/* /roms/ 2>/dev/null || true
    chown -R "${APP_USER:-abc}":"${APP_USER:-abc}" /roms/ 2>/dev/null || true
    echo "[svc-esde] Bundled demo ROMs copied to /roms. Add your own ROMs to replace them."
  fi
fi

# --- Pre-create ROM directory structure ---
# ES-DE shows a first-run "create directories?" wizard when the ROMs directory
# has no system subdirectories. After the user clicks "create", ES-DE loops back
# to the same wizard because there are still 0 games. Pre-creating the structure
# avoids this; ES-DE will instead show the normal empty-gamelist screen.
_roms_dir="${esde_home}/ROMs"
if [[ -d "${_roms_dir}" ]]; then
  _subdir_count="$(find "${_roms_dir}" -mindepth 1 -maxdepth 1 -type d 2>/dev/null | wc -l)"
  if [[ "${_subdir_count}" -eq 0 ]]; then
    echo "[svc-esde] ROM directory is empty; pre-creating system directories..."
    if s6-setuidgid "${APP_USER:-abc}" /usr/local/bin/es-de --create-system-dirs --home "${esde_home}" >/dev/null 2>&1; then
      echo "[svc-esde] System directories created successfully."
    else
      echo "[svc-esde] Warning: --create-system-dirs failed; ES-DE will prompt on first launch." >&2
    fi
  fi
fi

# --- Log forwarding to stdout (for balena / docker logs) ---
# ES-DE and Xorg write to files, not stdout. Tail them in the background
# so their output appears in `balena logs` / `docker logs`.
_log_pids=()
start_log_tailers() {
  local xorg_log="/var/log/Xorg.${display_num}.log"

  # Ensure the ES-DE log directory exists so tail -F can watch for creation
  mkdir -p "${esde_home}/logs"
  chown "${APP_USER:-abc}":"${APP_USER:-abc}" "${esde_home}/logs"

  # ES-DE may write its log to different paths depending on version and config.
  # Search known locations; fall back to the standard path for tail -F to watch.
  local esde_log=""
  local -a esde_log_candidates=(
    "${esde_home}/logs/es_log.txt"
    "${esde_home}/.emulationstation-de/es_log.txt"
    "${esde_home}/ES-DE/logs/es_log.txt"
    "${esde_home}/es_log.txt"
  )
  for _candidate in "${esde_log_candidates[@]}"; do
    if [[ -f "${_candidate}" ]]; then
      esde_log="${_candidate}"
      break
    fi
  done
  # If no existing log file found, watch the standard path for creation
  esde_log="${esde_log:-${esde_log_candidates[0]}}"
  echo "[svc-esde] Tailing ES-DE log: ${esde_log}"

  # Use -n 0 to skip existing content and only show new lines from this session
  tail -n 0 -F "${esde_log}" 2>/dev/null | sed --unbuffered 's/^/[es-de] /' &
  _log_pids+=("$!")
  tail -n 0 -F "${xorg_log}" 2>/dev/null | sed --unbuffered 's/^/[Xorg] /' &
  _log_pids+=("$!")
}

cleanup_log_tailers() {
  for pid in "${_log_pids[@]}"; do
    kill "${pid}" 2>/dev/null || true
  done
  # Also kill any tail children spawned in the pipelines
  jobs -p 2>/dev/null | xargs -r kill 2>/dev/null || true
}
trap cleanup_log_tailers EXIT

echo "[svc-esde] Starting ES-DE with DISPLAY=${DISPLAY}, XDG_RUNTIME_DIR=${XDG_RUNTIME_DIR}, and home=${esde_home}"
echo "[svc-esde] Launch target is local device graphics output; no web console or HTTP port is exposed by this image."

# --- udev setup for input device discovery (keyboard, mouse, gamepad) ---
# Non-balena base images do not auto-start udev; we must do it ourselves.
# Without a running udevd, Xorg cannot hotplug-discover input devices.
# In balena, io.balena.features.dbus mounts the socket at /run/dbus/system_bus_socket
# and sets DBUS_SYSTEM_BUS_ADDRESS automatically. However, if /run is a tmpfs the
# socket may be wiped. Verify the socket actually exists; start our own daemon if not.

_dbus_socket="/run/dbus/system_bus_socket"
_need_dbus=true

# _dbus_socket_alive: verify a D-Bus socket is actually responding, not just present on disk.
# A stale or bind-mounted socket file can exist without a listening daemon, causing Xorg to
# spam "dbus-core: error connecting to system bus" every 10 seconds.
_dbus_socket_alive() {
  local sock="$1"
  [[ -S "${sock}" ]] || return 1
  if command -v dbus-send >/dev/null 2>&1; then
    # Use GetId instead of Peer.Ping; the default system bus policy denies
    # Peer.Ping to the bus daemon itself, producing noisy "Rejected send" logs.
    DBUS_SYSTEM_BUS_ADDRESS="unix:path=${sock}" \
      dbus-send --system --dest=org.freedesktop.DBus --print-reply /org/freedesktop/DBus \
        org.freedesktop.DBus.GetId >/dev/null 2>&1
  else
    # No dbus-send available; fall back to file-type check only
    return 0
  fi
}

# Check if an existing dbus socket is reachable
if [[ -n "${DBUS_SYSTEM_BUS_ADDRESS:-}" ]]; then
  # Extract path from unix:path=... format
  _configured_socket="${DBUS_SYSTEM_BUS_ADDRESS#unix:path=}"
  if _dbus_socket_alive "${_configured_socket}"; then
    echo "[svc-esde] D-Bus system bus already available at ${_configured_socket}"
    _need_dbus=false
  else
    echo "[svc-esde] DBUS_SYSTEM_BUS_ADDRESS is set but socket ${_configured_socket} is not responding; will start our own." >&2
    unset DBUS_SYSTEM_BUS_ADDRESS
  fi
elif _dbus_socket_alive "${_dbus_socket}"; then
  export DBUS_SYSTEM_BUS_ADDRESS="unix:path=${_dbus_socket}"
  echo "[svc-esde] D-Bus system bus found at ${_dbus_socket}"
  _need_dbus=false
elif [[ -S "${_dbus_socket}" ]]; then
  echo "[svc-esde] D-Bus socket exists at ${_dbus_socket} but is not responding; will start our own." >&2
fi

if ${_need_dbus}; then
  if command -v dbus-daemon >/dev/null 2>&1; then
    if [[ ! -f /usr/share/dbus-1/system.conf ]]; then
      echo "[svc-esde] Warning: /usr/share/dbus-1/system.conf not found. Install 'dbus' package." >&2
      echo "[svc-esde]   Skipping dbus-daemon start; Xorg will log dbus-core errors (harmless but noisy)." >&2
    else
      if command -v dbus-uuidgen >/dev/null 2>&1; then
        dbus-uuidgen --ensure 2>/dev/null || true
      fi
      mkdir -p /run/dbus
      dbus-daemon --system --nofork --nopidfile 2>&1 | sed --unbuffered 's/^/[dbus] /' &
      # Wait for the socket to appear AND the daemon to respond to pings
      for _i in $(seq 1 50); do
        _dbus_socket_alive "${_dbus_socket}" && break
        sleep 0.1
      done
      if _dbus_socket_alive "${_dbus_socket}"; then
        export DBUS_SYSTEM_BUS_ADDRESS="unix:path=${_dbus_socket}"
        echo "[svc-esde] Started dbus-daemon; system bus available."
      else
        echo "[svc-esde] Warning: dbus-daemon started but not responding after 5s." >&2
      fi
    fi
  else
    echo "[svc-esde] Warning: dbus-daemon not installed. Xorg will log dbus-core errors." >&2
  fi
fi

# --- Video driver selection ---
# ESDE_VIDEO_DRIVER controls how ES-DE accesses the display:
#   kmsdrm  - Direct DRM/KMS (no X11). Best for SBCs like Raspberry Pi.
#   x11     - Start internal Xorg (legacy ESDE_USE_INTERNAL_X=1 path).
#   auto    - Auto-detect: prefer kmsdrm if /dev/dri/card0 exists.
# Falls through to external X11 socket if none of the above match.
_video_driver="${ESDE_VIDEO_DRIVER:-auto}"
if [[ "${_video_driver}" = "auto" ]]; then
  if [[ "${ESDE_USE_INTERNAL_X:-0}" = "1" ]]; then
    _video_driver="x11"
  elif [[ -e /dev/dri/card0 ]]; then
    _video_driver="kmsdrm"
  fi
fi
echo "[svc-esde] Video driver selection: ESDE_VIDEO_DRIVER=${ESDE_VIDEO_DRIVER:-<unset>}, resolved to '${_video_driver}'"

if [[ "${_video_driver}" = "kmsdrm" ]]; then
  echo "[svc-esde] Using KMSDRM video driver (direct to display, no Xorg)."
  export SDL_VIDEODRIVER=kmsdrm

  # Force SDL2 to use EGL/GLES instead of desktop GL.
  # Required for GPUs that only expose GLES (e.g. Raspberry Pi V3D, Mali).
  export SDL_OPENGL_ES_DRIVER=1

  # Prevent SDL2 from falling back to X11 or being confused by X11 env vars.
  unset DISPLAY 2>/dev/null || true
  unset XAUTHORITY 2>/dev/null || true

  # Log DRM devices for diagnostics
  if [[ -d /dev/dri ]]; then
    echo "[svc-esde] DRM devices: $(ls /dev/dri/ 2>/dev/null | tr '\n' ' ')"
  else
    echo "[svc-esde] Error: /dev/dri not found. KMSDRM requires GPU passthrough." >&2
  fi

  # KMSDRM needs direct VT access; ensure we have a usable VT.
  # SDL2 will open /dev/tty and issue VT_* ioctls.
  if [[ ! -e /dev/tty0 ]]; then
    echo "[svc-esde] Warning: /dev/tty0 not found. KMSDRM may fail without VT access." >&2
    echo "[svc-esde]   Ensure the container runs with --privileged or has /dev/tty* access." >&2
  fi

  # Raspberry Pi DRM architecture: card0 = vc4 display controller (GLES 2.0 only),
  # card1 = v3d GPU (GLES 3.1). When EGL creates a context on card0's GBM device
  # it may fail with EGL_BAD_MATCH because vc4 cannot provide the requested GLES
  # version. Detect this and fall back to X11 automatically.

  # KMSDRM requires DRM master (drmSetMaster ioctl), which needs root or
  # CAP_SYS_ADMIN. s6-setuidgid drops all capabilities, so es-de must run
  # as root for KMSDRM. This is acceptable for a dedicated appliance container.
  # The es-de home directory still needs to be writable.
  chown -R "${APP_USER:-abc}":"${APP_USER:-abc}" "${esde_home}" 2>/dev/null || true

  start_log_tailers
  echo "[svc-esde] Launching ES-DE via KMSDRM as root (DRM master requires root/CAP_SYS_ADMIN)."
  /usr/local/bin/es-de --home "${esde_home}" "${extra_args[@]}" 2>&1 | sed --unbuffered 's/^/[es-de] /'
  _kmsdrm_rc=${PIPESTATUS[0]}
  cleanup_log_tailers

  if [[ ${_kmsdrm_rc} -eq 0 ]]; then
    exit 0
  fi

  echo "[svc-esde] KMSDRM launch failed (exit code ${_kmsdrm_rc}). Falling back to X11 (Xorg)." >&2
  unset SDL_VIDEODRIVER SDL_OPENGL_ES_DRIVER
  export DISPLAY=":${display_num}"
  _video_driver="x11"
  # Fall through to X11 section below
fi

if [[ "${_video_driver}" = "x11" || "${ESDE_USE_INTERNAL_X:-0}" = "1" ]]; then
  echo "[svc-esde] Starting internal Xorg session via startx (as root for VT management)."

  # Clean stale X lock/socket files from a previous crash to prevent startx failure
  rm -f "/tmp/.X${display_num}-lock" 2>/dev/null || true
  rm -f "/tmp/.X11-unix/X${display_num}" 2>/dev/null || true

  sed -i -e 's/console/anybody/g' /etc/X11/Xwrapper.config || true
  if ! grep -q '^needs_root_rights=yes$' /etc/X11/Xwrapper.config 2>/dev/null; then
    echo 'needs_root_rights=yes' >> /etc/X11/Xwrapper.config
  fi
  export XAUTHORITY="${XAUTHORITY:-${XDG_RUNTIME_DIR}/.Xauthority}"
  touch "${XAUTHORITY}"
  chmod 600 "${XAUTHORITY}"
  chown "${APP_USER:-abc}":"${APP_USER:-abc}" "${XAUTHORITY}"
  if command -v mcookie >/dev/null 2>&1 && command -v xauth >/dev/null 2>&1; then
    cookie="$(mcookie)"
    xauth -f "${XAUTHORITY}" remove "${DISPLAY}" >/dev/null 2>&1 || true
    xauth -f "${XAUTHORITY}" add "${DISPLAY}" . "${cookie}" >/dev/null 2>&1 || true
  fi
  startx_server_args=("${DISPLAY}" -nolisten tcp -novtswitch -keeptty -nocursor -auth "${XAUTHORITY}")
  if [[ "${ESDE_XORG_VERBOSE:-0}" = "1" ]]; then
    startx_server_args+=( -verbose 3 )
  fi
  start_log_tailers
  # Xorg detailed output is already captured by the [Xorg] log tailer from the log file.
  # Filter startx output to only show non-Xorg messages (startx wrapper, xinit, errors).
  startx /usr/local/bin/esde-startx-session -- "${startx_server_args[@]}" 2>&1 \
    | grep -v --line-buffered -E '^\(([=-][=-]|[*][*]|\!\!|\?\?|II|WW|EE|NI)\)|^X\.Org X Server|^X Protocol Version|^Current Operating System|^Kernel command line|^xorg-server |^Current version of pixman|^[[:space:]]*(Before reporting|to make sure)|^Markers:|^[[:space:]]*\(' \
    | sed --unbuffered 's/^/[startx] /'
  exit ${PIPESTATUS[0]}
fi

if [[ ! -S /tmp/.X11-unix/X${DISPLAY#:} ]]; then
  echo "[svc-esde] X11 socket /tmp/.X11-unix/X${DISPLAY#:} is not available. Mount /tmp/.X11-unix from the host, set ESDE_VIDEO_DRIVER=x11, or set ESDE_VIDEO_DRIVER=kmsdrm." >&2
  sleep infinity
fi

start_log_tailers
if [[ -z ${LSIO_NON_ROOT_USER:-} ]]; then
  s6-setuidgid "${APP_USER}" /usr/local/bin/es-de --home "${esde_home}" "${extra_args[@]}" 2>&1 | sed --unbuffered 's/^/[es-de] /'
  exit ${PIPESTATUS[0]}
else
  /usr/local/bin/es-de --home "${esde_home}" "${extra_args[@]}" 2>&1 | sed --unbuffered 's/^/[es-de] /'
  exit ${PIPESTATUS[0]}
fi
